# triggercmd/rootfs/etc/services.d/triggercmd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Declare variables
declare token

## Get the 'message' key from the user config options.
token=$(bashio::config 'token')

## Print the message the user supplied, defaults to "Hello World..."
bashio::log.info "${token:="Missing"}"

echo "env in run:"
env

## Run your program
exec /usr/bin/my_program $token
